# docs/deploy/lgbapp-api.service.example
# Example systemd unit for the .NET API on Linux.
#   sudo cp docs/deploy/lgbapp-api.service.example /etc/systemd/system/lgbapp-api.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now lgbapp-api

[Unit]
Description=LGB Services API
After=network.target

[Service]
WorkingDirectory=/var/www/lgbapp/api
ExecStart=/usr/bin/dotnet /var/www/lgbapp/api/LGBApp.Backend.dll
Restart=always
RestartSec=5
User=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:5003
Environment=ConnectionStrings__DefaultConnection=Server=...;Database=LGBApp;...
Environment=Jwt__Key=...
Environment=Cors__AllowedOrigins__0=https://app.example.com

[Install]
WantedBy=multi-user.target
