# docs/deploy/lgbapp-api.service.example
# Example systemd unit for Lightsail UAT (Linux + RDS SQL Server).
#
#   sudo mkdir -p /etc/lgbapp
#   sudo cp docs/deploy/lightsail.env.example /etc/lgbapp/lgbapp.env
#   sudo nano /etc/lgbapp/lgbapp.env
#   sudo chmod 600 /etc/lgbapp/lgbapp.env
#   sudo cp docs/deploy/lgbapp-api.service.example /etc/systemd/system/lgbapp-api.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now lgbapp-api

[Unit]
Description=LGB Services API (UAT)
After=network.target

[Service]
WorkingDirectory=/var/www/lgbapp/api
ExecStart=/usr/bin/dotnet /var/www/lgbapp/api/LGBApp.Backend.dll
Restart=always
RestartSec=5
User=www-data
Group=www-data
EnvironmentFile=/etc/lgbapp/lgbapp.env

[Install]
WantedBy=multi-user.target
